# deploy/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ socketfile }}

[Install]
WantedBy=sockets.target
